# systemd/disk.service
[Unit]
Description=Disk Space Monitor

[Service]
Type=oneshot
ExecStart=/zpool/catallenya/ntfy/disk-ntfy.sh
