# restic/check/restic.check.timer
[Unit]
Description=Yearly Restic Check Timer

[Timer]
OnCalendar=*-06-01 3:00:00
Persistent=true

[Install]
WantedBy=timers.target